# systemd/user-negative-curfew-10.service
[Unit]
Description=Reminder: nightly shutdown (10 min lead)
ConditionPathExists=%h/code/futon0/scripts/negative_space_notify.py

[Service]
Type=oneshot
ExecStart=%h/code/futon0/scripts/negative_space_notify.py --title "Wrap up" --body "Final call: shut down by 21:00." --urgency normal
