# systemd/user-vitality-scanner.service
[Unit]
Description=Run FUTON vitality scanner snapshot
ConditionPathExists=%h/code/futon0/scripts/futon0/vitality/scanner.clj

[Service]
Type=oneshot
WorkingDirectory=%h/code/futon0
ExecStart=/usr/bin/env bb %h/code/futon0/scripts/futon0/vitality/scanner.bb -- --quiet --config %h/code/storage/futon0/vitality/vitality_scanner.json --output %h/code/futon3/resources/vitality/latest_scan.json
