# systemd/user-zoom-sync.path
[Unit]
Description=Trigger Zoom R4 ingest when the SD card is detected (user scope)
StartLimitIntervalSec=0
StartLimitBurst=0

[Path]
PathExistsGlob=/media/%u/R4_SD*/R4_Project
Unit=zoom-sync.service
TriggerLimitIntervalSec=0
TriggerLimitBurst=0

[Install]
WantedBy=default.target
