# systemd/user-negative-tai-chi-10.service
[Unit]
Description=Reminder: leave for Tai Chi (10 min lead)
ConditionPathExists=%h/code/futon0/scripts/negative_space_notify.py

[Service]
Type=oneshot
ExecStart=%h/code/futon0/scripts/negative_space_notify.py --title "Tai Chi" --body "Class at 19:00 → final call to leave now." --urgency normal
