# systemd/zoom-mount.path
[Unit]
Description=Detect Zoom R4 SD card by UUID

[Path]
PathExists=/dev/disk/by-uuid/6970-2CEF

[Install]
WantedBy=default.target
