# systemd/zoom-auto.target
[Unit]
Description=Run Zoom R4 ingest when the card appears
Requires=zoom-mount.path
