# systemd/zoom-mount.path
[Unit]
Description=Detect Zoom R4 SD card (user scope)

[Path]
PathChanged=/dev/disk/by-uuid

[Install]
WantedBy=default.target
